# systemd/ugreen-led-monitor.service
[Unit]
Description=UGREEN LED Monitor Service
After=network.target

[Service]
Type=simple
User=root
ExecStart=/opt/ugreen-led-controller/scripts/led_daemon.sh _daemon_process
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
